# templates/launchd-tunnel.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Tunnelclaw — autossh reverse tunnel plist. Generated by chrome-cdp-relay skill. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.chrome-tunnel-{{PROFILE_NAME}}</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/autossh</string>
    <!-- Disable autossh's built-in monitoring; rely on SSH keepalives instead -->
    <string>-M</string><string>0</string>
    <string>-N</string>
    <string>-i</string><string>{{SSH_KEY}}</string>
    <string>-o</string><string>ExitOnForwardFailure=yes</string>
    <string>-o</string><string>ServerAliveInterval=30</string>
    <string>-o</string><string>ServerAliveCountMax=3</string>
    <string>-o</string><string>StrictHostKeyChecking=accept-new</string>
    <!-- Reverse tunnel: remote:REMOTE_CDP_PORT → localhost:LOCAL_CDP_PORT -->
    <string>-R</string><string>127.0.0.1:{{REMOTE_CDP_PORT}}:127.0.0.1:{{LOCAL_CDP_PORT}}</string>
    <string>{{REMOTE_USER}}@{{REMOTE_HOST}}</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>{{HOME}}</string>
    <!-- AUTOSSH_GATETIME=0: don't require tunnel to be up for N seconds before giving up -->
    <key>AUTOSSH_GATETIME</key><string>0</string>
    <key>AUTOSSH_POLL</key><string>30</string>
    <key>AUTOSSH_FIRST_POLL</key><string>10</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>

  <key>StandardOutPath</key><string>/tmp/openclaw-autossh-{{PROFILE_NAME}}.out.log</string>
  <key>StandardErrorPath</key><string>/tmp/openclaw-autossh-{{PROFILE_NAME}}.err.log</string>
</dict>
</plist>
